# systemd/spi-button-controller.service
[Unit]
Description=SPI Button Controller Daemon
Documentation=man:spi-button-controller(8)
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/spi-button-controller /etc/spi-button-controller/config.yaml
ExecReload=/bin/kill -HUP $MAINPID

# User/Group - change these if needed
User=root
Group=root

# Restart policy
Restart=on-failure
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=spi-button-controller

# Security hardening (optional)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/dev/spidev*

[Install]
WantedBy=multi-user.target
